# live-save.service
[Unit]
Description=live-save
[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/live-save

[Install]
WantedBy=multi-user.target
