# debian/postinst
which systemctl && systemctl enable live-save.service
which rc-update && rc-update add live-save default
